# gain_ctrld-hack
#!/bin/sh /etc/rc.common
START=89
start(){
p=/data/system/wpa_supplicant.conf
cat /etc/wpa_supplicant.conf > $p 
cat >>$p<<EOF
network={
scan_ssid=1
ssid="{SSID}"
psk="{WIFI Pass}"
}
EOF
}
